# deploy/systemd/eimemory-karpathy-loop.service
[Unit]
Description=eimemory Karpathy Loop nightly cron (50 experiments / 4h budget)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/eimemory/current
Environment=EIMEMORY_ROOT=/var/lib/eimemory
Environment=EIMEMORY_CONFIG_DIR=/etc/eimemory
ExecStart=/opt/eimemory/current/scripts/karpathy_loop_cron.sh
# Wall-clock cap matches the in-script TIME_BUDGET_SECONDS=14400 (4h).
# We add a 5-minute headroom for the wrapper's own teardown so systemd
# does not SIGKILL the loop while it is writing its last kept-log row.
TimeoutStartSec=14700
# CPU cap so a runaway experiment cannot eat the whole box.
CPUQuota=200%
# Hard memory cap. The loop is mostly Python + JSONL; 4 GiB is generous
# and prevents a single bad hypothesis generator from OOM-ing the host.
MemoryMax=4G
