# deploy/systemd/eimemory-l5-observation-gate.service
[Unit]
Description=eimemory 48-hour observation gate for autonomous code changes
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/eimemory/current
Environment=EIMEMORY_ROOT=/var/lib/eimemory
Environment=EIMEMORY_CONFIG_DIR=/etc/eimemory
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONPYCACHEPREFIX=/var/lib/eimemory/.pycache/1.9.61
ExecStart=%h/.config/systemd/user/eimemory-l5-observation-gate.sh
