# deploy/systemd/eimemory-learn-dashboard.timer
[Unit]
Description=Run eimemory daily autonomous learning dashboard

[Timer]
OnCalendar=*-*-* 03:45:00
Persistent=true
Unit=eimemory-learn-dashboard.service

[Install]
WantedBy=timers.target
