# deploy/systemd/eimemory-console.service
[Unit]
Description=eimemory Evolution Console
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/lib/eimemory/governance
Environment=EIMEMORY_CONSOLE_TOKEN=change-me
Environment=EIMEMORY_CONSOLE_PATH=/var/lib/eimemory/governance/evolution-console.html
Environment=EIMEMORY_CONSOLE_PORT=8765
ExecStart=/opt/eimemory/venv/bin/python /var/lib/eimemory/governance/serve_console.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target

