# deploy/systemd/eimemory-learn-watch.service
[Unit]
Description=eimemory lightweight autonomous learning observer
After=network-online.target
Wants=network-online.target

[Service]
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONPYCACHEPREFIX=/var/lib/eimemory/.pycache/1.9.55
Type=oneshot
WorkingDirectory=/opt/eimemory/current
Environment=EIMEMORY_ROOT=/var/lib/eimemory
Environment=EIMEMORY_CONFIG_DIR=/etc/eimemory
ExecStart=/opt/eimemory/current/.venv/bin/eimemory learn watch --apply
