# deploy/systemd/eimemory-nightly.service
[Unit]
Description=eimemory Nightly Knowledge Intake and Governance
After=network-online.target
Wants=network-online.target

[Service]
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONPYCACHEPREFIX=/var/lib/eimemory/.pycache/1.9.59
Type=oneshot
WorkingDirectory=/opt/eimemory/current
Environment=EIMEMORY_ROOT=/var/lib/eimemory
Environment=EIMEMORY_CONFIG_DIR=/etc/eimemory
Environment=EIMEMORY_AUTONOMOUS_LEARNING_ENABLED=1
Environment=EIMEMORY_AUTONOMOUS_LEARNING_APPLY=0
Environment=EIMEMORY_AUTONOMOUS_LEARNING_DRY_RUN=0
Environment=EIMEMORY_AUTONOMOUS_LEARNING_MAX_GOALS=3
Environment=EIMEMORY_AUTONOMOUS_LEARNING_MAX_PROMOTIONS=3
Environment=EIMEMORY_AUTONOMOUS_LEARNING_NETWORK=1
Environment=EIMEMORY_AUTONOMOUS_LEARNING_TIMEOUT_SECONDS=900
Environment=EIMEMORY_AUTONOMOUS_LEARNING_DASHBOARD_ENABLED=1
Environment=EIMEMORY_L5_LOOP_ENABLED=1
Environment=EIMEMORY_L5_LOOP_APPLY=0
Environment=EIMEMORY_L5_LOOP_NETWORK=1
Environment=EIMEMORY_L5_MAX_GOALS=1
Environment=EIMEMORY_L5_MAX_PROMOTIONS=0
Environment=EIMEMORY_L5_LOOP_TIMEOUT_SECONDS=300
ExecStart=/opt/eimemory/current/.venv/bin/eimemory nightly
