# deploy/systemd/eimemory-learn-watch.timer
[Unit]
Description=Run eimemory lightweight observer every 15 minutes

[Timer]
OnCalendar=*:00/15
Persistent=true
Unit=eimemory-learn-watch.service

[Install]
WantedBy=timers.target
