# deploy/systemd/eimemory-karpathy-loop.timer
[Unit]
Description=Run the eimemory Karpathy Loop once per night (50 experiments / 4h)

[Timer]
# 02:30 every day in the server's local timezone. The 1.0.0 nightly
# intake already runs at 03:30; 02:30 leaves a 60-minute gap so the
# loop can finish (or be killed by the service's own time cap) before
# the nightly intake starts writing to the same state directory.
OnCalendar=*-*-* 02:30:00
# Run if the previous tick was missed (e.g. the host was off at 02:30).
Persistent=true
# The unit this timer activates.
Unit=eimemory-karpathy-loop.service

[Install]
WantedBy=timers.target
