# deploy/systemd/eimemory-rpc.service
[Unit]
Description=eimemory eibrain RPC service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/eimemory/current
Environment=EIMEMORY_ROOT=/var/lib/eimemory
Environment=EIMEMORY_CONFIG_DIR=/etc/eimemory
ExecStart=/opt/eimemory/current/.venv/bin/eimemory serve-eibrain-rpc --host 100.105.189.120 --port 8091
Restart=always
RestartSec=3
StandardOutput=append:%h/.openclaw/logs/eimemory-rpc.service.log
StandardError=append:%h/.openclaw/logs/eimemory-rpc.service.log

[Install]
WantedBy=default.target
