# deploy/systemd/openclaw-loop-watch.service
[Unit]
Description=OpenClaw loop watchdog
After=network-online.target openclaw-gateway.service eimemory-rpc.service

[Service]
Type=oneshot
WorkingDirectory=/opt/eimemory/current
Environment=OPENCLAW_LOOP_HOME=%h/.openclaw/task-ledger
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONPYCACHEPREFIX=/var/lib/eimemory/.pycache/1.9.26
ExecStart=/opt/eimemory/current/.venv/bin/python /opt/eimemory/current/scripts/openclaw_loop.py watch
