# deploy/systemd/openclaw-stuck-watchdog.service
[Unit]
Description=Restart OpenClaw gateway when sessions stall or hooks exhaust memory

[Service]
Type=oneshot
TimeoutStartSec=30
WorkingDirectory=/opt/eimemory/current
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONPYCACHEPREFIX=/var/lib/eimemory/.pycache/1.9.52
ExecStart=/opt/eimemory/current/.venv/bin/python -m eimemory.ops.openclaw_watchdog --unit openclaw-gateway.service --threshold-s 120 --min-restart-interval-s 300 --max-hook-processes 8 --max-hook-rss-mib 1536
